# sheng-sensors-files/usr/lib/systemd/system/iio-sensor-proxy.service.d/10-sheng-sensors.conf
[Unit]
Wants=adsprpcd-sensorspd.service
After=adsprpcd-sensorspd.service

[Service]
ExecStartPre=/bin/sleep 8
